# docker-web/config/supervisor.conf
[program:app]
command = uwsgi --ini /config/app.ini --catch-exceptions

[program:nginx]
command = service nginx restart